# football_bot.service
# /root/Complete_Seperate/football_bot.service
# Complete systemd unit file for Football Match Tracking System

[Unit]
# line 4-5: Basic service description and dependencies
Description=Football Match Tracking Bot
After=network.target

[Service]
# line 8-14: Service execution configuration
Type=simple
WorkingDirectory=/root/Complete_Seperate
ExecStart=/bin/bash /root/Complete_Seperate/run_pipeline.sh
Restart=on-failure
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1 LOG_STRICT=1

[Install]
# line 17: System target for auto-start
WantedBy=multi-user.target
