# football_bot_fixed.service
# /root/Complete_Seperate/football_bot_fixed.service
# line 1-2: Complete systemd unit file for Football Match Tracking System

[Unit]
# line 4-6: Service description and dependencies
Description=Football Match Tracking Bot (Fixed)
After=network.target
StartLimitIntervalSec=600

[Service]
# line 10-20: Service configuration with failsafe options
Type=simple
WorkingDirectory=/root/Complete_Seperate
# Use our robust wrapper script instead of run_pipeline.sh directly
ExecStart=/bin/bash /root/Complete_Seperate/service_wrapper.sh
# Restart only on clean exit and signals, not on abnormal exit
RestartSec=60
# Hard limit on runtime to prevent indefinite hanging
RuntimeMaxSec=900
User=root
# Set resource limits to prevent runaway processes
MemoryMax=256M
CPUQuota=50%
# Shorter timeouts to detect problems faster
TimeoutStartSec=120
TimeoutStopSec=30
# Essential environment variables
Environment=PYTHONUNBUFFERED=1 LOG_STRICT=0 API_TIMEOUT=30

[Install]
# line 23: System target for auto-start
WantedBy=multi-user.target
