# football_bot.timer
# /root/Complete_Seperate/football_bot.timer
# line 1-2: File header and description
# Timer unit to periodically execute the Football Match Tracking Bot

[Unit]
# line 5: Timer description
Description=Run Football Match Tracking Bot every 60 seconds

[Timer]
# lines 8-11: Timer configuration
# Wait 1 minute after boot before first run
OnBootSec=1min
# Schedule next run 60 seconds after previous service completion
OnUnitActiveSec=60s
# Ensure accurate timing within 1 second
AccuracySec=1s
# Reference to the service unit to trigger
Unit=football_bot_oneshot.service

[Install]
# line 16: System target for auto-start
WantedBy=timers.target
