# football_bot_oneshot.service
# /root/Complete_Seperate/football_bot_oneshot.service
# line 1-2: File header and description
# Complete systemd unit file for Football Match Tracking System (oneshot execution)

[Unit]
# lines 5-7: Service description and dependencies
Description=Football Match Tracking Bot (oneshot)
After=network.target
Documentation=https://github.com/alextrx818/Complete_Seperate

[Service]
# lines 10-16: Service execution configuration for oneshot operation
Type=oneshot
WorkingDirectory=/root/Complete_Seperate
# line 13: Correct path to virtualenv for this specific deployment
ExecStart=/bin/bash -lc 'source /root/Complete_Seperate/sports_venv/bin/activate && exec /root/Complete_Seperate/run_pipeline.sh'
User=root
TimeoutStartSec=300
Environment=PYTHONUNBUFFERED=1 LOG_STRICT=0 API_TIMEOUT=60
ExecStartPre=/bin/bash -c 'echo "$(date) ENV: $VIRTUAL_ENV" >> /root/Complete_Seperate/logs/env_debug.log'
